# deploy/supervisor.conf
[program:unblockme]
command=twistd -n unblockme -i 199.195.251.154 -d .netlix.com -c home.kraiz.de,146-52-6-45-dynip.superkabel.de
directory=/home/kraiz/unblockme
autostart=true
autorestart=true