# docs/old_setup/proxy.pac/supervisor.conf
[program:proxy]
command=gunicorn -w1 proxy.py
directory=/home/kraiz/
user=kraiz
autostart=true
autorestart=true
stdout_logfile=/home/kraiz/proxy.out.log
stderr_logfile=/home/kraiz/proxy.err.log